# deploy/diag-ingest.timer
[Unit]
Description=Weekly RAG Corpus Ingest — runs every Sunday at 2am UTC

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
